# etc-init-uwsgi.conf
# simple uWSGI script

description "uWSGI - Beer"
start on runlevel [2345]
stop on runlevel [06]

exec sudo uwsgi --die-on-term --yaml /etc/uwsgi/beer.yaml
